# systemd/nfs-config.service
[Unit]
Description=Preprocess NFS configuration

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/systemd/scripts/nfs-utils_env.sh
