# systemd/proc-fs-nfsd.mount
[Unit]
Description=NFSD configuration filesystem
DefaultDependencies=no

[Mount]
What=nfsd
Where=/proc/fs/nfsd
Type=nfsd
