# debian/nfs-common.postrm
#!/bin/sh

#DEBHELPER#

case "$1" in
    purge)
	update-rc.d nfs-common remove >/dev/null

	rm -f /var/lib/nfs/state    \
	      /var/lib/nfs/sm/*     \
	      /var/lib/nfs/sm.bak/*
	;;
esac
