# debian/nfs-common.postinst
#!/bin/sh -e

#DEBHELPER#

case "$1" in
    configure)
	update-rc.d nfs-common defaults 19 81 >/dev/null
	;;
esac

act="restart"
[ "$1:$2" = "configure:" ] && act="start"
/etc/init.d/nfs-common $act
