# debian/nfs-kernel-server.postrm
#!/bin/sh

#DEBHELPER#

case "$1" in
    purge)
	update-rc.d nfs-kernel-server remove >/dev/null
	;;
esac

