# debian/nfs-common.conffiles
/etc/default/nfs-common
/etc/init.d/nfs-common
