# debian/nfs-kernel-server.postrm
#!/bin/sh

#DEBHELPER#

case "$1" in
    purge)
	update-rc.d nfs-kernel-server remove >/dev/null

	rm -f /var/lib/nfs/etab  \
              /var/lib/nfs/rmtab \
              /var/lib/nfs/xtab
	;;
esac

